# services/simpleTopologyService/supervisor_config/mongo.conf
[program:mongo]
command=/usr/bin/mongod --dbpath /data/db
