# services/simpleTopologyService/supervisor_config/node.conf
[program:node]
command=node /app/app.js
